# roles/mojibake/files/mojibake.ini
[program:mojibake]
process_name=mojibake%(process_num)s
directory=/opt/mojibake
environment=USER=mojibake,PYTHONPATH=/opt/mojibake/venv/bin
environment=LANG="en_AU.UTF-8",LC_ALL="en_AU.UTF-8"
command=/opt/mojibake/venv/bin/python /opt/mojibake/tornado_srv.py
startsecs=2
user=mojibake
autostart=true
autorestart=true
numprocs=1
numprocs_start=8080
stderr_logfile = /var/log/supervisor/mojibake-stderr.log
stdout_logfile = /var/log/supervisor/mojibake-stdout.log
